# fan-lock.service
[Unit]
Description=Keep fan at a fixed speed (see fan-lock.sh)
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/fan-lock.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
